# templates/ubuntu/lstacker
#!/bin/bash
### BEGIN INIT INFO
# Short-Description: lstacker
# Description: The startup scripts provided for Lstacker
# This is auto generated by lstacker during install
### END INIT INFO

LSTACKER_DIR=/etc/lstacker

case "$1" in
 start)
   ${LSTACKER_DIR}/startup.sh
   ;;
 stop)
   ${LSTACKER_DIR}/startup.sh
   ;;
 *)
   echo "Usage: lstacker {start|stop}" >&2
   exit 3
   ;;
esac
